# systemd/claude-monitor.service
[Unit]
Description=Claude Usage Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yoxie
WorkingDirectory=/home/yoxie/claude_progress
ExecStart=/usr/bin/python3 /home/yoxie/claude_progress/claude_lcd.py
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
